# crontab/root
# At every 10th minute past every hour from 7 through 23.
*/10 7-23,0 * * * node /app/dist/index.js > /dev/stdout